# config/th-mon.timer
[Unit]
Description=Run th-mon every five minutes

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
